# .config/systemd/user/irexec.service
[Unit]
Description=Irexec
Requires=lirc.service
After=lirc.service

[Service]
ExecStart=/usr/bin/irexec

[Install]
WantedBy=mystuff.target
