# .config/systemd/user/mbsync.service
[Unit]
Description=Synchronize IMAP4 and Maildir mailboxes

[Service]
ExecStart=/usr/bin/mbsync -a

[Install]
WantedBy=timer-30min.target
