# .config/systemd/user/timer-30min.target
[Unit]
Description=30 Min Timer Target
StopWhenUnneeded=yes
